# etc/periodic/daily/Makefile
# $FreeBSD: src/etc/periodic/daily/Makefile,v 1.4.2.4 2002/07/18 12:36:07 ru Exp $

FILES=	100.clean-disks \
	110.clean-tmps \
	120.clean-preserve \
	130.clean-msgs \
	140.clean-rwho \
	150.clean-hoststat \
	160.clean-hammer \
	161.clean-hammer2 \
	200.backup-passwd \
	210.backup-aliases \
	300.calendar \
	310.accounting \
	330.news \
	400.status-disks \
	420.status-network \
	430.status-uptime \
	440.status-mailq \
	450.status-security \
	460.status-mail-rejects \
	500.queuerun \
	999.local

.include <bsd.prog.mk>
