# etc/rc.d/addswap
#!/bin/sh
#
# Add additional swap files
#
# $FreeBSD: src/etc/rc.d/addswap,v 1.2 2002/10/12 10:31:31 schweikh Exp $
#

# PROVIDE: addswap
# REQUIRE: mountcritlocal
# BEFORE: sysctl

. /etc/rc.subr

name="addswap"
start_cmd="addswap_start"
stop_cmd=":"

addswap_start()
{
	case ${swapfile} in
	[Nn][Oo] | '')
		;;
	*)
		if [ -w "${swapfile}" ]; then
			# load vn kernel module if needed
			kldload -n vn || exit 1
			echo "Adding ${swapfile} as additional swap"
			vnconfig -e vn0 ${swapfile} swap
		fi
		;;
	esac
}

load_rc_config $name
run_rc_command "$1"
