# etc/rc.d/devfs
#!/bin/sh
#
# Load default devfs(5) rules.
#

# PROVIDE: devfs
# REQUIRE: rcconf
# BEFORE: disks initrandom

. /etc/rc.subr

name="devfs"

load_rc_config $name
dummy_rc_command "$1"

load_rules()
{
	echo -n "Loading devfs rules:"
	set -- $devfs_config_files
	for file; do
		[ -r $file ] && echo -n " $file" && /sbin/devfsctl -a -f $file
	done
	echo '.'
}

case "$1" in
*start)
	load_rules
	;;
esac
