# etc/rc.d/ldconfig
#!/bin/sh
#
# $NetBSD: ldconfig,v 1.5 2002/03/22 04:33:58 thorpej Exp $
# $FreeBSD: src/etc/rc.d/ldconfig,v 1.7 2003/06/30 15:02:05 trhodes Exp $
# $DragonFly: src/etc/rc.d/ldconfig,v 1.5 2004/03/20 16:27:39 drhodus Exp $
#

# PROVIDE: ldconfig
# REQUIRE: mountall mountcritremote
# BEFORE:  DAEMON
# KEYWORD: DragonFly

. /etc/rc.subr

name="ldconfig"
ldconfig_command="/sbin/ldconfig"
start_cmd="ldconfig_start"
stop_cmd=":"

ldconfig_start()
{
		_ins=
		ldconfig=${ldconfig_command}
		checkyesno ldconfig_insecure && _ins="-i"
		if [ -x "${ldconfig_command}" ]; then
			_LDC="/usr/lib /usr/lib/gcc*"
			for i in ${ldconfig_paths} /etc/ld-elf.so.conf; do
				if [ -r "${i}" ]; then
					_LDC="${_LDC} ${i}"
				fi
			done
			echo 'ELF ldconfig path:' ${_LDC}
			${ldconfig} -elf ${_ins} ${_LDC}
		fi
}

load_rc_config $name
run_rc_command "$1"
