# etc/rc.d/pppoed
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/pppoed,v 1.1 2002/06/13 22:14:36 gordon Exp $
# $DragonFly: src/etc/rc.d/pppoed,v 1.2 2004/01/27 00:42:45 rob Exp $
#

# PROVIDE: pppoed
# REQUIRE: NETWORKING
# BEFORE: DAEMON
# KEYWORD: DragonFly

. /etc/rc.subr

name="pppoed"
rcvar="`set_rcvar`"
start_cmd="pppoed_start"
# XXX stop_cmd will not be straightforward
stop_cmd=":"

pppoed_start()
{
	if [ -n "${pppoed_provider}" ]; then
			pppoed_flags="${pppoed_flags} -p ${pppoed_provider}"
	fi
	echo 'Starting pppoed'
	_opts=$-; set -f
	/usr/libexec/pppoed ${pppoed_flags} ${pppoed_interface}
	set +f; set -${_opts}
}

load_rc_config $name
run_rc_command "$1"
