# etc/rc.d/rcconf
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/rcconf.sh,v 1.2 2003/01/25 20:02:35 mtm Exp $
# $DragonFly: src/etc/rc.d/rcconf.sh,v 1.4 2005/11/19 21:47:32 swildner Exp $
#

# PROVIDE: rcconf
# REQUIRE: initdiskless
# BEFORE:  devfs

. /etc/rc.subr
dummy_rc_command "$1"

echo "Loading configuration files."
load_rc_config 'XXX'
