# etc/rc.d/sshd
#!/bin/sh
#
# $NetBSD: sshd,v 1.18 2002/04/29 08:23:34 lukem Exp $
# $FreeBSD: src/etc/rc.d/sshd,v 1.3 2003/07/13 01:49:07 mtm Exp $
# $DragonFly: src/etc/rc.d/sshd,v 1.3 2005/11/19 21:47:32 swildner Exp $
#

# PROVIDE: sshd
# REQUIRE: LOGIN

. /etc/rc.subr

name="sshd"
rcvar=`set_rcvar`
keygen_cmd="sshd_keygen"
start_precmd="sshd_precmd"
pidfile="/var/run/${name}.pid"
extra_commands="keygen reload"

sshd_keygen()
{
	/usr/bin/ssh-keygen -A
}

sshd_precmd()
{
	run_rc_command keygen
}

load_rc_config $name
run_rc_command "$1"
